# install/mac/com.jcu.auto-network.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 服务标识符 -->
    <key>Label</key>
    <string>com.jcu.auto-network</string>
    
    <!-- 程序路径和参数 -->
    <key>ProgramArguments</key>
    <array>
        <string>PYTHON_PATH_PLACEHOLDER</string>
        <string>APP_PATH_PLACEHOLDER</string>
    </array>
    
    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <string>WORKING_DIR_PLACEHOLDER</string>
    
    <!-- 自动启动配置 -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- 保持运行 - 如果程序退出会自动重启 -->
    <key>KeepAlive</key>
    <dict>
        <!-- 网络可用时才启动 -->
        <key>NetworkState</key>
        <true/>
        <!-- 如果意外退出，自动重启 -->
        <key>Crashed</key>
        <true/>
    </dict>
    
    <!-- 仅捕获错误输出，程序自身已有详细日志 -->
    <!-- StandardOutPath 删除，避免与程序日志重复 -->
    
    <key>StandardErrorPath</key>
    <string>WORKING_DIR_PLACEHOLDER/logs/jcu-auto-network-error.log</string>
    
    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>PYTHONPATH_PLACEHOLDER</string>
    </dict>
    
    <!-- 启动延迟（等待系统完全启动）-->
    <key>StartInterval</key>
    <integer>30</integer>
    
    <!-- 资源限制 -->
    <key>SoftResourceLimits</key>
    <dict>
        <!-- 最大内存使用 (256MB) -->
        <key>MemoryLimit</key>
        <integer>268435456</integer>
    </dict>
    
    <!-- 进程类型 -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- 退出超时 -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
    
    <!-- 限制启动频率，避免频繁重启 -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>